# package/etc/init/redis-server.conf
description "redis server"

expect daemon

start on runlevel [23]
stop on shutdown

exec sudo -u redis /usr/bin/redis-server /etc/redis/redis.conf

respawn
